# config/systemd/mempalace-chroma-server.service
[Unit]
Description=MemPalace shared ChromaDB HTTP daemon
Documentation=https://github.com/crewrig/crewrig/blob/main/docs/adr/0006-chromadb-http-server.md
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.mempalace
ExecStart=%h/.local/pipx/venvs/mempalace/bin/python3.13 %h/.local/pipx/venvs/mempalace/bin/chroma run --path %h/.mempalace/palace --host 127.0.0.1 --port 8001
Restart=always
RestartSec=3
StandardOutput=append:%h/.mempalace/chroma-server.log
StandardError=append:%h/.mempalace/chroma-server.log

[Install]
WantedBy=default.target
